# install/autonomy-netmon.service
[Unit]
Description=Autonomy Network Monitor Daemon
Documentation=https://github.com/Autonomy-Logic/orchestrator-agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/autonomy-netmon
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autonomy-netmon

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/orchestrator /var/log

[Install]
WantedBy=multi-user.target
